# systemd/cnc-wifi-fallback.service
[Unit]
Description=CNC Control Wi-Fi fallback
After=NetworkManager.service network.target
Wants=NetworkManager.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/cnc-control/cnc-control.env
ExecStart=/bin/bash -lc 'exec "${CNC_CONTROL_REPO:-/home/cnc/cnc-control}/tools/wifi_fallback.sh"'

[Install]
WantedBy=multi-user.target
